# ansible/system-config/roles/os_role/files/Redhat/profile.d/lingering.sh
if (( ${EUID} >= 1000 )) ; then
  export XDG_RUNTIME_DIR=/run/user/$(id -u)
fi
